# Library/profile.d/asdf.sh
if [ -f "$HOME/.asdf/asdf.sh" ]; then
  . $HOME/.asdf/asdf.sh;
fi
