# Library/profile.d/catalina-deprecation.sh
# removes the annoying "The default interactive shell is now zsh."
export BASH_SILENCE_DEPRECATION_WARNING=1
