# Library/profile.d/clearf.sh
#clear iTerm scrollback history
alias clearf="printf '\e]50;ClearScrollback\a'";
