# Library/profile.d/dotlocal.sh
# enable bin from $HOME/.local
if [ -d "$HOME/.local/bin" ]; then
    export PATH="$PATH:$HOME/.local/bin"
fi