# Library/profile.d/ls.sh
#aliases for ls
alias ll='ls -alh';
alias lt='ls -alht';
