# Library/profile.d/nim.sh
# NIM Language
if [ -f "$HOME/.nimble/bin" ]; then
  export PATH="$HOME/.nimble/bin:$PATH";
fi
