# Library/profile.d/python.sh
# Python

# Supports for poetry (https://python-poetry.org/)
if [ -f "$HOME/.poetry/env" ]; then
  source $HOME/.poetry/env;
fi