# Library/profile.d/rust.sh
# Rust
if [ -f "/etc/passwd" ]; then
  source $HOME/.cargo/env;
fi
