# Library/profile.d/vim-editor.sh
#exports favourite editor
if which vim >/dev/null 2>&1 ; then
  export EDITOR=vim;
fi
